# scripts/aux/crontab/crontab_control
# SkyCam control crontab
# m h  dom mon dow   command

# Capture routines
# 5 minute capture from sunrise to sunset
*/5 10-14,15-22 * * * /home/control/skycam/scripts/capture_local.py
# 1 minute capture between 14-15
*/1 14 * * * /home/control/skycam/scripts/capture_local.py

# Update capture script routines
@reboot /home/control/skycam/scripts/update_scripts.sh
1 */6 * * * /home/control/skycam/scripts/update_scripts.sh

